# Deploy/systemd/sea-of-fire-backend.service
[Unit]
Description=Sea of fire backend
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/sea-of-fire/backend
EnvironmentFile=/opt/sea-of-fire/backend/.env
ExecStart=/usr/bin/npm start
Restart=always
RestartSec=5
User=ubuntu

[Install]
WantedBy=multi-user.target

